# env_logger.service


[Unit]
Description=Log temperature etc
After=multi-user.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/pi-environment-logger/pi_environment_logger.py
User=pi

[Install]
WantedBy=multi-user.target


# Copy content of this file to:
# sudo systemctl edit --force --full  env_logger.service
# Set environment variables at:
# sudo nano /etc/systemd/system/env_logger.service.d/override.conf
# sudo systemctl daemon-reload
# sudo systemctl start env_logger.service